# ansible/roles/grafana/templates/etc/systemd/system/alertmanager.service.j2
[Unit]
Description=Alertmanager handles alerts generated by prometheus and other systems
Requires=network.target
After=network.target
StartLimitIntervalSec=0

[Service]
WorkingDirectory=/home/{{ alertmanager_system_user }}
User={{ alertmanager_system_user }}
Group={{ alertmanager_system_user }}
ExecStart=/home/{{ alertmanager_system_user }}/bin/alertmanager --config.file=/home/{{ alertmanager_system_user }}/conf/alertmanager.yml --cluster.peer={{ alertmanager_cluster_peers[ansible_ec2_instance_identity_document_availabilityzone] }} --storage.path /home/{{ alertmanager_system_user}}/data
ExecStop=/home/{{ alertmanager_system_user }}/bin/alertmanager --remove-pidfile
StandardOutput=append:/var/log/{{ alertmanager_systemd_service_name }}.log
StandardError=append:/var/log/{{ alertmanager_systemd_service_name }}.log
KillMode=mixed
TimeoutSec=60
Type=simple
#Environment=LOG_DIR="/home/{{ alertmanager_system_user }}/logs/"
LimitNOFILE=100000
Restart=always
RestartSec=5s

[Install]
WantedBy=multi-user.target
